# supervisor-metube.conf
[program:metube-uwsgi]
command = /usr/local/bin/uwsgi --ini /usr/local/metube/uwsgi.ini

[program:metube-nginx]
command = /usr/sbin/nginx
